# debian/qt-autoconnect.cron.d
MAILTO=email@example.com
@reboot qt-autoconnect >/dev/null 2>&1
0 4 * * * qt-autoconnect >/dev/null 2>&1
0 10 * * * qt-autoconnect >/dev/null 2>&1
0 16 * * * qt-autoconnect >/dev/null 2>&1
0 22 * * * qt-autoconnect >/dev/null 2>&1
